# etc/service/dnsmasq-safe/run
#!/bin/sh

exec [ -e /etc/dnsmasq.d/safe ] && /usr/sbin/dnsmasq -k -x /var/run/dnsmasq/dnsmasq-safe.pid -p 5353 -u dnsmasq -7 /etc/dnsmasq.d/safe,.dpkg-dist,.dpkg-old,.dpkg-new || true
